# badcat.service
[Unit]
Description=Input Lock
Documentation=man:badcat(8)

[Service]
ExecStart=badcat

[Install]
WantedBy=multi-user.target
